# Project/crontab
30 16 * * * root python3 /app/financeInsert.py >> /var/log/cron.log 2>&1